# overlay/system/etc/init/fpc-key-disable.conf
description "FPC Key Disabler"

start on started lightdm

respawn

exec /usr/bin/uinput-fpc-key-disable
