# overlay/system/etc/init/goodix-key-disable.conf
description "Goodix Key Disabler"

start on started lightdm

respawn

exec /usr/bin/uinput-goodix-key-disable
