# overlay/system/etc/init/usb-tethering.conf
# upstart service file for usb-tethering

start on startup

exec /bin/bash /usr/bin/usb-tethering
